# monitor-agent/tutor-monitor-agent.timer
[Unit]
Description=Run Tutor IA monitor agent every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=tutor-monitor-agent.service

[Install]
WantedBy=timers.target
